# .config/killchrome.cron
#!/bin/bash
/usr/bin/ps axo pid,pcpu,comm | /usr/bin/grep chrome | /usr/bin/awk '{ if($2 > 90) print $1 }' | /usr/bin/xargs kill > /dev/null 2>&1

